# test/resources/language/systemd-analyze-test.service
[Unit]
Description=Test file for flycheck systemd-analyze checker
Documentation=foo://bar

[Service]
ExecSmart=

[Dog]
